# arch-ganglia/build/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
loglevel=WARN
logfile_maxbytes=0
username=root

[program:gmetad]
command=/usr/bin/gmetad
username=root
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

[program:gmond]
command=/usr/bin/gmond
username=root
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

[program:php-fm]
command=usr/bin/php-fpm
username=root
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

[program:nginx]
command=/usr/bin/nginx
username=root
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
